# files/etc/uci-defaults/98-set-root-password
#!/bin/sh
# 这个脚本在首次启动时运行，用于设置 root 密码
# 脚本执行成功后会被自动删除

# 设置 root 密码为 12345678
# 使用 passwd 命令的非交互方式
echo -e "12345678\n12345678" | passwd root

exit 0
